# btrbk-backup.service
[Unit]
Description=Run BTRBK backup to external drive
Conflicts=shutdown.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=oneshot
EnvironmentFile=-/etc/linuxscripts/backup.conf
# Wait up to 30 seconds for the mountpoint to appear
ExecStartPre=/usr/bin/timeout 30s bash -c 'until mountpoint -q "${BACKUP_MOUNT:-/path/to/your/backup/mount}"; do sleep 1; done'
ExecStart=/usr/local/bin/backup-notify.sh
User=root
StandardOutput=journal

[Install]
WantedBy=multi-user.target
