# btrbk-backup.service
[Unit]
Description=Run BTRBK backup to external drive
Conflicts=shutdown.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=oneshot
# Wait up to 30 seconds for the mountpoint to appear
# Replace with your specific mount path
ExecStartPre=/usr/bin/timeout 30s bash -c 'until mountpoint -q /path/to/your/backup/mount; do sleep 1; done'
ExecStart=/usr/local/bin/backup-notify.sh
User=root
StandardOutput=journal

[Install]
WantedBy=multi-user.target
